# app/redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server
stdout_events_enabled=true
stderr_events_enabled=true